# packages/res-scheduling/res-scheduling-1.0.2.tar.gz/res-scheduling-1.0.2/res/scheduling/res_scheduling.service
[Unit]
Description=RES Scheduling Service
After=rabbitmq-server.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/usr/lib/python3/dist-packages
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=res.scheduling
ExecStart=/usr/share/res/scheduling/run.sh
Restart=always
RestartSec=1
WatchdogSec=10
OOMScoreAdjust=-1000

[Install]
WantedBy=multi-user.target
